# s6/services.d/deluge/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid deluge /usr/bin/deluged -d -c /config/deluged --loglevel=info -l /config/deluged/deluged.log
